# deploy/systemd/bot-metrics-web.service
[Unit]
Description=bot-metrics — Flask dashboard (gunicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bot
Group=bot
WorkingDirectory=/opt/bot-metrics/repo
EnvironmentFile=/opt/bot-metrics/.env
ExecStart=/opt/bot-metrics/venv/bin/gunicorn \
  --bind 127.0.0.1:8765 \
  --workers 2 \
  --timeout 30 \
  --access-logfile - \
  --error-logfile - \
  web:app
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
